# rtkbase/s6-rc.d/rtk-caster/run
#!/command/with-contenv bashio

exec /usr/local/rtkbase/rtkbase/run_cast.sh
